# misc/nftables/nftables-restore.timer
[Unit]
Description=Wait to restore /etc/nftables.conf

[Timer]
OnActiveSec=5m
RemainAfterElapse=false
