# misc/NetworkManager/dispatcher.d/99-chrony
#!/bin/bash

IFACE="$1"
ACTION="$2"

exec /etc/admin/scripts/misc/chrony.sh "$ACTION"
